# etc/systemd/my-tracks.service
[Unit]
Description=My Tracks – OwnTracks location tracking server
Documentation=https://github.com/the-hcma/my-tracks
After=network.target

[Service]
Type=simple
WorkingDirectory=@@REPO_DIR@@
ExecStart=@@REPO_DIR@@/scripts/my-tracks-server --http-port 8080 --log-level debug --console --log-file %h/scratch/my-tracks/my-tracks.log
ExecStartPost=/bin/bash -c 'i=0; while [ $i -lt 30 ]; do sleep 1; curl -sf http://localhost:8080/api/health/ && exit 0; i=$((i+1)); done; exit 1'
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
